# infrastructure/terraform-staging/files/genbi-backend.service
[Unit]
Description=GenBI Backend API
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=genbi
Group=genbi
WorkingDirectory=/opt/genbi/backend
EnvironmentFile=/opt/genbi/backend/.env
ExecStart=/opt/genbi/backend/venv/bin/uvicorn src.main:app --host 127.0.0.1 --port 8000 --workers 2
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
